# files/pg_back.cron
30 23 * * *   postgres    /opt/pg_back/pg_back -q -c /opt/pg_back/pg_back.conf
